# deploy/shivad.service
[Unit]
Description=Shiva blockchain node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shiva
Group=shiva
WorkingDirectory=/opt/shiva
ExecStart=/opt/shiva/bin/shivad \
  -datadir /var/lib/shiva \
  -genesis /opt/shiva/configs/genesis.json \
  -seeds /opt/shiva/configs/seeds-mainnet.json \
  -api :8545 \
  -listen :30303 \
  -tls-cert /etc/shiva/tls/fullchain.pem \
  -tls-key /etc/shiva/tls/privkey.pem
Environment=SHIVA_CORS_ORIGINS=https://your-domain.com
EnvironmentFile=-/etc/shiva/shiva.env
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
